# rootfs/etc/services.d/pi1/run
#!/usr/bin/with-contenv sh
/go/bin/mjpeg-server -a :9000 -b FRAME -- curl -N $PI1_URL