# rootfs/etc/services.d/pi2/run
#!/usr/bin/with-contenv sh
/go/bin/mjpeg-server -a :9001 -b FRAME -p "/stream.mjpg" -- curl -sN $PI2_URL