# rootfs/etc/services.d/pi2/run
#!/usr/bin/with-contenv sh
/go/bin/mjpeg-server -a :9001 -b FRAME -- curl -N $PI2_URL